# bin/service/linkguardian-prod.service
[Unit]
Description=Linkguardian production
[Service]
User=linkguardian
# The configuration file application.properties should be here:
#change this to your workspace
WorkingDirectory=/home/linkguardian/apps
#path to executable.
#executable is a bash script which calls jar
ExecStart=/home/linkguardian/apps/start-linkguardian-production.sh
SuccessExitStatus=143
TimeoutStopSec=10
Restart=on-failure
RestartSec=5
[Install]
WantedBy=multi-user.target
